# luci-app-terranet-dmz-pro/root/etc/init.d/terranet_dmz_pro
#!/bin/sh /etc/rc.common
START=99

start()  { /usr/lib/terranet_dmz_pro.sh; }
stop()   { /usr/lib/terranet_dmz_pro.sh; }
